# time.cron
*/10 * * * * /bin/bash -C '/home/ecs-user/project/HotSearchDataAnalysis/exec.sh'
*/60 * * * * /bin/bash -C '/home/ecs-user/project/HotSearchDataAnalysis/exec-spark.sh'