# time.cron
*/2 * * * * /bin/bash -C ~/project/HotSearchDataAnalysis/exec.sh
